# Codigo Raspberry/pi/.config/autostart/crawler2.desktop
[Desktop Entry]
Type=Application
Name=Crawler_Server2
NoDisplay=true
Exec=touch /home/pi/dola.txt && python3 app.py > /home/pi/error_app.txt

